# supervisor-mailhog.conf

[program:mailhog]
command=/usr/local/bin/mailhog
stopsignal=6
